# packaging/debian/prerm
#!/bin/bash
set -e

if command -v systemctl >/dev/null 2>&1; then
    systemctl stop eshield.service || true
    systemctl disable eshield.service || true
fi

exit 0
